# code/ansible/roles/bmc_simulator/templates/bmc-sol-simulator.service.j2
# File: bmc-sol-simulator.service.j2
[Unit]
Description=BMC Serial Over LAN (SOL) Simulator
After=network.target bmc-core-simulator.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/bmc-simulator
Environment="PYTHONUNBUFFERED=1"
ExecStart=/opt/bmc-simulator/venv/bin/python /opt/bmc-simulator/bmc_sol_simulator.py
Restart=always
RestartSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target